# ubuntu-rootfs/etc/systemd/system/snap-lxd-31333.mount
[Unit]
Description=Mount unit for lxd, revision 31333
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/lxd_31333.snap
Where=/snap/lxd/31333
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
