# ubuntu-rootfs/etc/systemd/system/snap-snapd-23545.mount
[Unit]
Description=Mount unit for snapd, revision 23545
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/snapd_23545.snap
Where=/snap/snapd/23545
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
