# ubuntu-rootfs/var/lib/dpkg/info/finalrd.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/finalrd.service
/usr
/usr/bin
/usr/bin/finalrd
/usr/lib
/usr/lib/finalrd
/usr/lib/finalrd/finalrd-static.conf
/usr/share
/usr/share/doc
/usr/share/doc/finalrd
/usr/share/doc/finalrd/changelog.gz
/usr/share/doc/finalrd/copyright
/usr/share/finalrd
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/finalrd.1.gz
