# ubuntu-rootfs/var/lib/dpkg/info/lxd-agent-loader.list
/.
/lib
/lib/systemd
/lib/systemd/lxd-agent-setup
/lib/systemd/system
/lib/systemd/system/lxd-agent.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/99-lxd-agent.rules
/usr
/usr/share
/usr/share/doc
/usr/share/doc/lxd-agent-loader
/usr/share/doc/lxd-agent-loader/changelog.gz
/usr/share/doc/lxd-agent-loader/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/lxd-agent-loader
