# ubuntu-rootfs/var/lib/dpkg/info/multipath-tools.list
/.
/etc
/etc/multipath.conf
/lib
/lib/libdmmp.so.0.2.0
/lib/libmpathcmd.so.0
/lib/libmpathpersist.so.0
/lib/libmultipath.so.0
/lib/multipath
/lib/multipath/libcheckcciss_tur.so
/lib/multipath/libcheckdirectio.so
/lib/multipath/libcheckemc_clariion.so
/lib/multipath/libcheckhp_sw.so
/lib/multipath/libcheckrdac.so
/lib/multipath/libcheckreadsector0.so
/lib/multipath/libchecktur.so
/lib/multipath/libforeign-nvme.so
/lib/multipath/libprioalua.so
/lib/multipath/libprioana.so
/lib/multipath/libprioconst.so
/lib/multipath/libpriodatacore.so
/lib/multipath/libprioemc.so
/lib/multipath/libpriohds.so
/lib/multipath/libpriohp_sw.so
/lib/multipath/libprioiet.so
/lib/multipath/libprioontap.so
/lib/multipath/libpriopath_latency.so
/lib/multipath/libpriorandom.so
/lib/multipath/libpriordac.so
/lib/multipath/libpriosysfs.so
/lib/multipath/libprioweightedpath.so
/lib/pkgconfig
/lib/pkgconfig/libdmmp.pc
/lib/systemd
/lib/systemd/system
/lib/systemd/system/multipathd.service
/lib/systemd/system/multipathd.socket
/lib/tmpfiles.d
/lib/tmpfiles.d/multipath.conf
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/56-dm-mpath.rules
/lib/udev/rules.d/60-multipath.rules
/sbin
/sbin/mpathpersist
/sbin/multipath
/sbin/multipathd
/usr
/usr/include
/usr/include/libdmmp
/usr/include/libdmmp/libdmmp.h
/usr/include/mpath_cmd.h
/usr/include/mpath_persist.h
/usr/share
/usr/share/bug
/usr/share/bug/multipath-tools
/usr/share/bug/multipath-tools/script
/usr/share/doc
/usr/share/doc/multipath-tools
/usr/share/doc/multipath-tools/README.Debian.gz
/usr/share/doc/multipath-tools/README.alua
/usr/share/doc/multipath-tools/README.md
/usr/share/doc/multipath-tools/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/multipath-tools
/usr/share/man
/usr/share/man/man3
/usr/share/man/man3/dmmp_context_free.3.gz
/usr/share/man/man3/dmmp_context_log_func_set.3.gz
/usr/share/man/man3/dmmp_context_log_priority_get.3.gz
/usr/share/man/man3/dmmp_context_log_priority_set.3.gz
/usr/share/man/man3/dmmp_context_new.3.gz
/usr/share/man/man3/dmmp_context_timeout_get.3.gz
/usr/share/man/man3/dmmp_context_timeout_set.3.gz
/usr/share/man/man3/dmmp_context_userdata_get.3.gz
/usr/share/man/man3/dmmp_context_userdata_set.3.gz
/usr/share/man/man3/dmmp_flush_mpath.3.gz
/usr/share/man/man3/dmmp_last_error_msg.3.gz
/usr/share/man/man3/dmmp_log_priority_str.3.gz
/usr/share/man/man3/dmmp_mpath_array_free.3.gz
/usr/share/man/man3/dmmp_mpath_array_get.3.gz
/usr/share/man/man3/dmmp_mpath_kdev_name_get.3.gz
/usr/share/man/man3/dmmp_mpath_name_get.3.gz
/usr/share/man/man3/dmmp_mpath_wwid_get.3.gz
/usr/share/man/man3/dmmp_path_array_get.3.gz
/usr/share/man/man3/dmmp_path_blk_name_get.3.gz
/usr/share/man/man3/dmmp_path_group_array_get.3.gz
/usr/share/man/man3/dmmp_path_group_id_get.3.gz
/usr/share/man/man3/dmmp_path_group_priority_get.3.gz
/usr/share/man/man3/dmmp_path_group_selector_get.3.gz
/usr/share/man/man3/dmmp_path_group_status_get.3.gz
/usr/share/man/man3/dmmp_path_group_status_str.3.gz
/usr/share/man/man3/dmmp_path_status_get.3.gz
/usr/share/man/man3/dmmp_path_status_str.3.gz
/usr/share/man/man3/dmmp_reconfig.3.gz
/usr/share/man/man3/dmmp_strerror.3.gz
/usr/share/man/man3/libdmmp.h.3.gz
/usr/share/man/man3/mpath_persistent_reserve_in.3.gz
/usr/share/man/man3/mpath_persistent_reserve_out.3.gz
/usr/share/man/man5
/usr/share/man/man5/multipath.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/mpathpersist.8.gz
/usr/share/man/man8/multipath.8.gz
/usr/share/man/man8/multipathd.8.gz
/lib/libdmmp.so
/lib/libmpathcmd.so
/lib/libmpathpersist.so
/lib/libmultipath.so
/lib/systemd/system/multipath-tools-boot.service
/lib/systemd/system/multipath-tools.service
/usr/share/doc/multipath-tools/NEWS.Debian.gz
/usr/share/doc/multipath-tools/changelog.Debian.gz
