# ubuntu-rootfs/var/lib/dpkg/info/xfsprogs.list
/.
/lib
/lib/libhandle.so.1.0.3
/lib/systemd
/lib/systemd/system
/lib/systemd/system/xfs_scrub@.service
/lib/systemd/system/xfs_scrub_all.service
/lib/systemd/system/xfs_scrub_all.timer
/lib/systemd/system/xfs_scrub_fail@.service
/sbin
/sbin/fsck.xfs
/sbin/mkfs.xfs
/sbin/xfs_repair
/usr
/usr/include
/usr/include/xfs
/usr/lib
/usr/lib/xfsprogs
/usr/lib/xfsprogs/xfs_scrub_all.cron
/usr/lib/xfsprogs/xfs_scrub_fail
/usr/sbin
/usr/sbin/xfs_admin
/usr/sbin/xfs_bmap
/usr/sbin/xfs_copy
/usr/sbin/xfs_db
/usr/sbin/xfs_estimate
/usr/sbin/xfs_freeze
/usr/sbin/xfs_fsr
/usr/sbin/xfs_growfs
/usr/sbin/xfs_info
/usr/sbin/xfs_io
/usr/sbin/xfs_logprint
/usr/sbin/xfs_mdrestore
/usr/sbin/xfs_metadump
/usr/sbin/xfs_mkfile
/usr/sbin/xfs_ncheck
/usr/sbin/xfs_quota
/usr/sbin/xfs_rtcp
/usr/sbin/xfs_scrub
/usr/sbin/xfs_scrub_all
/usr/sbin/xfs_spaceman
/usr/share
/usr/share/doc
/usr/share/doc/xfsprogs
/usr/share/doc/xfsprogs/CREDITS
/usr/share/doc/xfsprogs/README
/usr/share/doc/xfsprogs/changelog.Debian.gz
/usr/share/doc/xfsprogs/copyright
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/xfs
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/projects.5.gz
/usr/share/man/man5/projid.5.gz
/usr/share/man/man5/xfs.5.gz
/usr/share/man/man8
/usr/share/man/man8/fsck.xfs.8.gz
/usr/share/man/man8/mkfs.xfs.8.gz
/usr/share/man/man8/xfs_admin.8.gz
/usr/share/man/man8/xfs_bmap.8.gz
/usr/share/man/man8/xfs_copy.8.gz
/usr/share/man/man8/xfs_db.8.gz
/usr/share/man/man8/xfs_estimate.8.gz
/usr/share/man/man8/xfs_freeze.8.gz
/usr/share/man/man8/xfs_fsr.8.gz
/usr/share/man/man8/xfs_growfs.8.gz
/usr/share/man/man8/xfs_info.8.gz
/usr/share/man/man8/xfs_io.8.gz
/usr/share/man/man8/xfs_logprint.8.gz
/usr/share/man/man8/xfs_mdrestore.8.gz
/usr/share/man/man8/xfs_metadump.8.gz
/usr/share/man/man8/xfs_mkfile.8.gz
/usr/share/man/man8/xfs_ncheck.8.gz
/usr/share/man/man8/xfs_quota.8.gz
/usr/share/man/man8/xfs_repair.8.gz
/usr/share/man/man8/xfs_rtcp.8.gz
/usr/share/man/man8/xfs_scrub.8.gz
/usr/share/man/man8/xfs_scrub_all.8.gz
/usr/share/man/man8/xfs_spaceman.8.gz
/lib/libhandle.so.1
/usr/share/doc/xfsprogs/changelog.gz
